# pwnomcp.service
[Unit]
Description=Pwno MCP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pwno
WorkingDirectory=/app
Environment=PYTHONPATH=/app
Environment=UV_PROJECT_ENVIRONMENT=/app/.venv
ExecStart=/app/.venv/bin/python -m pwnomcp --host 0.0.0.0 --port 5500
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
